# adapters/openclaw/monitor.service
[Unit]
Description=Clade Monitor
After=network.target

[Service]
Type=simple
User=%i
Environment=CLADE_API_KEY=changeme
ExecStart=/usr/bin/python3 %h/projects/clade/adapters/openclaw/monitor.py --project %h/projects/myproject --port 9100
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
